# webdav-server/rootfs/etc/services.d/init/run
#!/usr/bin/with-contenv bash

# Exécuter le script d'initialisation
exec /usr/local/bin/init.sh
