# webdav-server/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/command/with-contenv sh
echo "Initializing WebDAV server..."
exec /usr/local/bin/init.sh

