# webdav-server/rootfs/etc/s6-overlay/s6-rc.d/apache/run
#!/command/with-contenv sh
echo "Starting Apache server..."
exec /usr/sbin/httpd -D FOREGROUND

