# webdav-server/rootfs/etc/s6-overlay/s6-rc.d/apache/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Apache server..."
exec apache2ctl -D FOREGROUND
